# Docker/Solr/supervisord.conf
[supervisord]
nodaemon=true

[program:solr]
command=/etc/init.d/solr start
process_name=%(program_name)s
startsecs=5
stopsignal=INT
redirect_stderr=true
stdout_logfile=/var/solr/logs/solr.log